# rootfs/etc/s6-overlay/s6-rc.d/rspamd/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
foreground { rspamadm configtest }
rspamd -f -u _rspamd -g _rspamd
