# rootfs/etc/s6-overlay/s6-rc.d/dovecot/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
foreground { doveconf -n }
dovecot -F
