# rootfs/etc/s6-overlay/s6-rc.d/postfix/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
foreground { postfix check }
postfix start-fg
