# rootfs/etc/s6-overlay/s6-rc.d/postfix/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
foreground { postfix check }
# Tee postfix stdout into the audit FIFO (when audit is enabled) or /dev/null.
# Using execlineb pipeline keeps postfix as s6's direct supervised child so
# s6 restarts the real postfix master on exit, not a shell wrapper.
#
# The shell one-liner resolves AUDIT_ENABLED and prints the tee target path so
# execlineb can capture it via backtick; the FIFO is created here so the tee
# open below never races against audit-svc's own mkfifo.
backtick -n AUDIT_FIFO {
  /bin/sh -c "
case \"\${AUDIT_ENABLED:-true}\" in
  0|false|FALSE|False|no|off) printf /dev/null; exit ;;
esac
mkdir -p /run/audit
[ -p /run/audit/maillog.fifo ] || mkfifo -m 0600 /run/audit/maillog.fifo
printf /run/audit/maillog.fifo"
}
importas AUDIT_FIFO AUDIT_FIFO
pipeline { postfix start-fg }
tee -a $AUDIT_FIFO
