# rootfs/etc/s6-overlay/s6-rc.d/postsrsd/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
importas -D "" MAIL_HOSTNAME MAIL_HOSTNAME
postsrsd -s /etc/postsrsd.secret -d ${MAIL_HOSTNAME} -f 10001 -r 10002 -l 127.0.0.1 -4
