# rootfs/etc/s6-overlay/s6-rc.d/mail-bootstrap/up
#!/command/execlineb -P
foreground { /usr/local/bin/mail-bootstrap }
