# rootfs/etc/s6-overlay/s6-rc.d/audit-svc/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
python3 -u /usr/local/bin/audit-svc.py
